# swadiq-app.service
[Unit]
Description=Swadiq Schools Application
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-/usr/bin/docker stop swadiq-app
ExecStartPre=-/usr/bin/docker rm swadiq-app
ExecStartPre=-/usr/bin/docker pull swadiq-schools:latest
ExecStart=/usr/bin/docker run -d --name swadiq-app -p 8080:8080 --restart unless-stopped swadiq-schools:latest
ExecStop=/usr/bin/docker stop swadiq-app

[Install]
WantedBy=multi-user.target